# server/deploy/server-monitor.service
[Unit]
Description=Server Monitor
After=network-online.target
Wants=network-online.target

[Service]
User=servermon
Group=servermon
WorkingDirectory=/opt/server-monitor
ExecStart=/opt/server-monitor/server -config /etc/server-monitor/config.yaml
Restart=on-failure
RestartSec=5
Environment=SERVER_MONITOR_LISTEN=:8443
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target



